# vpn.service
[Unit]
Description=vpn

[Service]
ExecStart=/usr/bin/vpn.sh

[Install]
WantedBy=multi-user.target
